# roles/common/files/mackerel-agent.service
[Unit]
Description=Run mackerel agent

[Service]
Type=simple
ExecStart=/usr/local/mackerel-agent/mackerel-agent --conf=/usr/local/mackerel-agent/mackerel-agent.conf
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
Restart=always

[Install]
WantedBy=multi-user.target
RequiredBy=multi-user.target
